# debian/ti-omap3-sgx.powervr.init
#!/sbin/openrc-run
# Copyright 2019 Maemo Leste
# Distributed under the terms of the GNU General Public License v2

depend() {
    before logger
    after eudev
}

start() {
    ebegin "Starting powervr"
    /usr/sbin/pvrsrvinit
}

stop() {
    ebegin "Stopping powervr (no-op)"
}
